# Hardware_Dev/ServiceScripts/schedule_logs.service
[Unit]
Description=Schedule Logs TempHum Readings

[Service]
ExecStart=templogs.sh

[Install]
WantedBy=multi-user.target
