# Hardware_Dev/ServiceScripts/schedule_logs.timer
[Unit]
Description=Run Logs Script every 1 hour

[Timer]
OnBootSec=3min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
